# debian/DEBIAN/postinst
#!/bin/bash
set -e

# Создаем пользователя mssql, если не существует (для монтирования)
if ! id "mssql" &>/dev/null; then
    useradd -r -s /bin/false mssql
fi

# Создаем директории, если не существуют
mkdir -p /var/log/sqlmanager
mkdir -p /mnt/sql_backups

# Устанавливаем права на директории и файлы
chown -R mssql:mssql /opt/SQLManager
chown -R mssql:mssql /var/log/sqlmanager

# Устанавливаем права: директории с rwx, файлы с rw
find /opt/SQLManager -type d -exec chmod 755 {} \;
find /opt/SQLManager -type f -exec chmod 644 {} \;

# Делаем исполняемым основной бинарник
chmod +x /opt/SQLManager/sqlmanager

# Устанавливаем специальные права для конфигурационного файла
if [ -f /opt/SQLManager/config.yaml ]; then
    chmod 600 /opt/SQLManager/config.yaml
fi

# Создаем, если отсутствует, файл в /etc/smbcredentials/.veeamsrv_creds для выполнения команд монтирования
if [ -f /etc/smbcredentials/.veeamsrv_creds ]; then
    echo "username=user"  > /etc/smbcredentials/.veeamsrv_creds
    echo "password=password"  > /etc/smbcredentials/.veeamsrv_creds
    echo "domain=domain" > /etc/smbcredentials/.veeamsrv_creds
    chmod 640 /etc/smbcredentials/.veeamsrv_creds
fi

# Устанавливаем права на лог-файл, если он существует
if [ -f /var/log/sqlmanager/sqlmanager.log ]; then
    chmod 640 /var/log/sqlmanager/sqlmanager.log
fi
chmod 750 /var/log/sqlmanager

# Создаем, если отсутствует, файл в /etc/sudoers.d для выполнения команд монтирования
if [ -f /etc/sudoers.d/sqlmanager ]; then
    echo "mssql ALL=(ALL) NOPASSWD: /bin/systemctl start mnt-sql_backups.mount, /bin/systemctl status mnt-sql_backups.mount" > /etc/sudoers.d/sqlmanager
    chmod 440 /etc/sudoers.d/sqlmanager
fi

# Перезагружаем systemd
systemctl daemon-reload || true

# Включаем и запускаем сервисы
systemctl enable mnt-sql_backups.mount || true
# systemctl start mnt-sql_backups.mount || true
systemctl enable sqlmanager.service || true
# systemctl start sqlmanager.service || true

echo
echo "SQLManager установлен."
echo "!!! Пожалуйста, отредактируйте файлы: /etc/smbcredentials/.veeamsrv_creds и /opt/SQLManager/config.yaml, если устанавливаете приложение впервые !!!"
echo "После настройки запустите службы: sudo systemctl start mnt-sql_backups.mount и sudo systemctl start sqlmanager.service или перезагрузите сервер."
