# debian/DEBIAN/postinst
#!/bin/bash
set -e

# Создаем пользователя mssql, если не существует (для монтирования)
if ! id "mssql" &>/dev/null; then
    useradd -r -s /bin/false mssql
fi

# Создаем файл в /etc/sudoers.d для выполнения команд монтирования
echo "mssql ALL=(ALL) NOPASSWD: /bin/systemctl start mnt-sql_backups.mount, /bin/systemctl status mnt-sql_backups.mount" > /etc/sudoers.d/sqlmanager
chmod 440 /etc/sudoers.d/sqlmanager

# Создаем директории, если не существуют
mkdir -p /var/log/sqlmanager
mkdir -p /mnt/sql_backups
mkdir -p /etc/smbcredentials

# Устанавливаем права на директории и файлы
chown -R mssql:mssql /opt/SQLManager
chown -R mssql:mssql /var/log/sqlmanager

# Устанавливаем права: директории с rwx, файлы с rw
find /opt/SQLManager -type d -exec chmod 755 {} \;
find /opt/SQLManager -type f -exec chmod 644 {} \;

# Делаем исполняемым основной бинарник
chmod +x /opt/SQLManager/sqlmanager

# Устанавливаем специальные права для конфигурационного файла
chmod 600 /opt/SQLManager/config.yaml

# Устанавливаем права на лог-файл, если он существует
if [ -f /var/log/sqlmanager/sqlmanager.log ]; then
    chmod 640 /var/log/sqlmanager/sqlmanager.log
fi
chmod 750 /var/log/sqlmanager

# Копируем файлы сервисов, если не существуют
if [ ! -f /etc/systemd/system/mnt-sql_backups.mount ]; then
    cp /opt/SQLManager/mnt-sql_backups.mount /etc/systemd/system/
    systemctl enable mnt-sql_backups.mount || true
fi

if [ ! -f /etc/systemd/system/sqlmanager.service ]; then
    cp /opt/SQLManager/sqlmanager.service /etc/systemd/system/
    systemctl enable sqlmanager.service || true
fi

# Перезагружаем systemd
systemctl daemon-reload || true

# Включаем и запускаем сервисы
systemctl enable mnt-sql_backups.mount || true
systemctl start mnt-sql_backups.mount || true
systemctl enable sqlmanager.service || true
systemctl start sqlmanager.service || true

echo "SQLManager установлен. Пожалуйста, настройте /opt/SQLManager/config.yaml и /etc/smbcredentials/.veeamsrv_creds перед использованием."
